# systemd/mpd-automount.service.in
[Unit]
Description=Mount volumes and integrate them into the MPD library.
After=mpd.service

[Service]
ExecStart=@MPD_AUTOMOUNT_BINARY@
Restart=always

[Install]
WantedBy=default.target

